# demostracion/poky/build/meta-ai/recipes-ai/email-agent/files/email-agent.service
[Unit]
Description=Email Sales Agent - Asistente de ventas por correo
After=network-online.target ollama.service
Wants=network-online.target
Requires=ollama.service
[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/bin/email-agent/agent.py
User=root
Environment=HOME=/root
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
[Install]
WantedBy=multi-user.target
